# systemd/quant-unified-monitor.service
[Unit]
Description=🦞 龙虾王量化系统 - 统一监控器
After=network.target

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/home/admin/.openclaw/workspace/quant/v3-architecture
ExecStart=/home/admin/.openclaw/workspace/quant/.venv/bin/python -u scripts/unified_monitor.py
Restart=always
RestartSec=10
StartLimitBurst=10
StartLimitInterval=60
Environment=PYTHONPATH=/home/admin/.openclaw/workspace/quant/v3-architecture
StandardOutput=append:/home/admin/.openclaw/workspace/quant/v3-architecture/logs/systemd_unified_monitor.log
StandardError=append:/home/admin/.openclaw/workspace/quant/v3-architecture/logs/systemd_unified_monitor_err.log

[Install]
WantedBy=multi-user.target
