# supervisor/auto_strategies.conf
# Auto-generated strategy configurations
# DO NOT EDIT MANUALLY

[program:quant-strategy-rsi_1min_strategy]
command=/root/.pyenv/versions/3.10.13/bin/python3 -u strategies/rsi_1min_strategy.py
directory=/root/.openclaw/workspace/quant/v3-architecture
user=root
group=root
autostart=true
autorestart=true
startretries=10
retrywait=5
stderr_logfile=/root/.openclaw/workspace/quant/v3-architecture/logs/supervisor_rsi_1min_strategy_err.log
stdout_logfile=/root/.openclaw/workspace/quant/v3-architecture/logs/supervisor_rsi_1min_strategy_out.log
environment=PYTHONPATH="/root/.openclaw/workspace/quant/v3-architecture"
numprocs=1
process_name=%(program_name)s
startsecs=5
stopwaitsecs=10
stopsignal=TERM

[program:quant-strategy-link_rsi_detailed_strategy]
command=/root/.pyenv/versions/3.10.13/bin/python3 -u strategies/link_rsi_detailed_strategy.py
directory=/root/.openclaw/workspace/quant/v3-architecture
user=root
group=root
autostart=true
autorestart=true
startretries=10
retrywait=5
stderr_logfile=/root/.openclaw/workspace/quant/v3-architecture/logs/supervisor_link_rsi_detailed_strategy_err.log
stdout_logfile=/root/.openclaw/workspace/quant/v3-architecture/logs/supervisor_link_rsi_detailed_strategy_out.log
environment=PYTHONPATH="/root/.openclaw/workspace/quant/v3-architecture"
numprocs=1
process_name=%(program_name)s
startsecs=5
stopwaitsecs=10
stopsignal=TERM

[program:quant-strategy-rsi_scale_in_strategy]
command=/root/.pyenv/versions/3.10.13/bin/python3 -u strategies/rsi_scale_in_strategy.py
directory=/root/.openclaw/workspace/quant/v3-architecture
user=root
group=root
autostart=true
autorestart=true
startretries=10
retrywait=5
stderr_logfile=/root/.openclaw/workspace/quant/v3-architecture/logs/supervisor_rsi_scale_in_strategy_err.log
stdout_logfile=/root/.openclaw/workspace/quant/v3-architecture/logs/supervisor_rsi_scale_in_strategy_out.log
environment=PYTHONPATH="/root/.openclaw/workspace/quant/v3-architecture"
numprocs=1
process_name=%(program_name)s
startsecs=5
stopwaitsecs=10
stopsignal=TERM

