# supervisor/quant-enhanced-monitor.conf
[program:quant-enhanced-monitor]
command=/home/admin/.pyenv/versions/3.10.0/bin/python3 -u scripts/enhanced_monitor.py
directory=/home/admin/.openclaw/workspace/quant/v3-architecture
user=admin
group=admin
autostart=true
autorestart=true
startretries=10
retrywait=5
stderr_logfile=/home/admin/.openclaw/workspace/quant/v3-architecture/logs/supervisor_enhanced_monitor_err.log
stdout_logfile=/home/admin/.openclaw/workspace/quant/v3-architecture/logs/supervisor_enhanced_monitor_out.log
environment=PYTHONPATH="/home/admin/.openclaw/workspace/quant/v3-architecture"
numprocs=1
process_name=%(program_name)s
startsecs=5
stopwaitsecs=10
stopsignal=TERM
