# supervisor/quant-strategies.conf
# V3 策略 Supervisor 配置（清理版）
# 只保留 3 个核心策略

[program:quant-strategy-eth]
command=/root/.pyenv/versions/3.10.13/bin/python3 -u strategies/rsi_1min_strategy.py
directory=/root/.openclaw/workspace/quant/v3-architecture
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=10
stderr_logfile=logs/supervisor_eth_err.log
stdout_logfile=logs/supervisor_eth_out.log
environment=PYTHONPATH="/root/.openclaw/workspace/quant/v3-architecture"

[program:quant-strategy-link]
command=/root/.pyenv/versions/3.10.13/bin/python3 -u strategies/link_rsi_detailed_strategy.py
directory=/root/.openclaw/workspace/quant/v3-architecture
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=10
stderr_logfile=logs/supervisor_link_err.log
stdout_logfile=logs/supervisor_link_out.log
environment=PYTHONPATH="/root/.openclaw/workspace/quant/v3-architecture"

[program:quant-strategy-avax]
command=/root/.pyenv/versions/3.10.13/bin/python3 -u strategies/rsi_scale_in_strategy.py
directory=/root/.openclaw/workspace/quant/v3-architecture
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=10
stderr_logfile=logs/supervisor_avax_err.log
stdout_logfile=logs/supervisor_avax_out.log
environment=PYTHONPATH="/root/.openclaw/workspace/quant/v3-architecture"
