# supervisor/quant-web.conf
[program:quant-web]
command=/home/admin/.pyenv/versions/3.10.0/bin/python3 -m uvicorn web.dashboard_api:app --host 0.0.0.0 --port 3000
directory=/home/admin/.openclaw/workspace/quant/v3-architecture
user=admin
group=admin
autostart=true
autorestart=true
startretries=10
retrywait=5
stderr_logfile=/home/admin/.openclaw/workspace/quant/v3-architecture/logs/supervisor_web_err.log
stdout_logfile=/home/admin/.openclaw/workspace/quant/v3-architecture/logs/supervisor_web_out.log
environment=PYTHONPATH="/home/admin/.openclaw/workspace/quant/v3-architecture"
numprocs=1
process_name=%(program_name)s
startsecs=5
stopwaitsecs=10
stopsignal=TERM
