# supervisor/quant-deep-monitor.conf
[program:quant-deep-monitor]
command=python3 scripts/deep_monitor.py
directory=/home/admin/.openclaw/workspace/quant/v3-architecture
autostart=true
autorestart=true
startretries=3
stderr_logfile=/home/admin/.openclaw/workspace/quant/v3-architecture/logs/deep_monitor_err.log
stdout_logfile=/home/admin/.openclaw/workspace/quant/v3-architecture/logs/deep_monitor_out.log
environment=PYTHONPATH="/home/admin/.openclaw/workspace/quant/v3-architecture"
