# supervisor/web_dashboard.conf
# 🦞 大王量化 Web Dashboard 守护进程配置

[program:web_dashboard]
command=/home/admin/.pyenv/versions/3.10.0/bin/python3 -m uvicorn web.dashboard_api:app --host 0.0.0.0 --port 3000
directory=/home/admin/.openclaw/workspace/quant/v3-architecture
user=admin
autostart=true
autorestart=true
startretries=3
stderr_logfile=/home/admin/.openclaw/workspace/quant/v3-architecture/logs/web_dashboard_err.log
stdout_logfile=/home/admin/.openclaw/workspace/quant/v3-architecture/logs/web_dashboard_out.log
environment=PYTHONPATH="/home/admin/.openclaw/workspace/quant/v3-architecture"
numprocs=1
process_name=%(program_name)s_%(process_num)02d
startsecs=5
stopwaitsecs=10
stopsignal=TERM

# 自动重启条件
# - 意外退出时自动重启
# - 最多重试 3 次
# - 启动后 5 秒内退出视为启动失败
