# supervisor/web_dashboard.conf
[program:web_dashboard]
command=/root/.pyenv/versions/3.10.13/bin/python3 -m uvicorn web.dashboard_api:app --host 0.0.0.0 --port 3000 --timeout-keep-alive 30
directory=/root/.openclaw/workspace/quant/v3-architecture
user=root
autostart=true
autorestart=true
startretries=3
stderr_logfile=/root/.openclaw/workspace/quant/v3-architecture/logs/web_dashboard_err.log
stdout_logfile=/root/.openclaw/workspace/quant/v3-architecture/logs/web_dashboard_out.log
environment=PYTHONPATH="/root/.openclaw/workspace/quant/v3-architecture"
startsecs=5
stopwaitsecs=30
stopsignal=TERM
killasgroup=true
stopasgroup=true
